# node-deb/upstart.conf
# This file was autogenerated by node-deb {{ node_deb_version }}

description "{{ node_deb_package_name }}"
author "{{ node_deb_package_maintainer }}"

start on runlevel [2345]
stop on runlevel [016]
respawn

# set the working directory of the job processes
chdir '{{ node_deb_install_dir }}/{{ node_deb_package_name }}/app'

# start the process
script
  . "/etc/default/{{ node_deb_package_name }}"
  exec sudo -u '{{ node_deb_user }}' /usr/local/bin/node app/{{ daemon_entrypoint }}
end script