# setup/OctoFarm.DaemonService/OctoFarmDaemon.service
[Unit]
Description=OctoFarm DaemonService

[Service]
Type=notify
ExecStart=/usr/sbin/octofarm-daemon

[Install]
WantedBy=multi-user.target